# py/flask_learning/myweb/supervisor_app.conf
[include]
files=/etc/supervisor/supervisord.conf

[program:myweb]
directory=/var/www/back_end/myweb
command=gunicorn -c gunicorn.py main:app
